# minecraft
#!/bin/sh
# /etc/init.d/minecraft
# version 0.5 2012-01-31 (YYYY-MM-DD)

### BEGIN INIT INFO
# Provides:   minecraft
# Required-Start: $local_fs $remote_fs
# Required-Stop:  $local_fs $remote_fs
# Should-Start:   $network
# Should-Stop:    $network
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# Short-Description:    Minecraft server
# Description:    Init script for bukkit server.
### END INIT INFO

# Based on http://www.minecraftwiki.net/wiki/Server_startup_script

#############################
######### SETTINGS ##########
#############################

# Name of server.jar file
SERVICE="craftbukkit-1.1.0-R3.jar"

# Name to use for the screen instance
SCREEN="minecraft110"

# User that should run the server
USERNAME="minecraft"

# Path to minecraft directory
MCBASEPATH="/home/${USERNAME}"

# Initial memory usage
INITMEM="512M"

# Maximum amount of memory to use
MAXMEM="1024M"









#############################
####### DO NOT EDIT #########
#############################

INVOCATION="java -server -Xmx$MAXMEM -Xms$INITMEM -jar $SERVICE $OPTIONS"

# Where the whole minecraft directory is copied when whole-backup is executed
MCPATH="${MCPATH}/server"

# Where the whole minecraft directory is copied when whole-backup is executed
BACKUP="${MCPATH}/server.bak"


screen_running(){
	if [ `screen -ls $SCREEN_NAME | grep $SCREEN_NAME` == "" ]; then
		return false;
	fi
	return true;
}

ME=`whoami`
as_user() {
	if [ $ME == $USERNAME ]; then
		sh -c "$1"
	else
		su - $USERNAME -c "$1"
	fi
}

mc_start() {
	if [ screen_running ]; then
		echo "Tried to start but $SERVICE was already running!"
	else
		echo "$SERVICE was not running... starting."
		cd $MCPATH
		as_user "cd $MCPATH && screen -dmS $SCREEN $INVOCATION"
		sleep 3
		if screen_running
		then
			echo "$SERVICE is now running."
		else
			echo "Could not start $SERVICE."
		fi
	fi
}

mc_saveoff() {
	if [ screen_running ]; then
		echo "$SERVICE is running... suspending saves"
		as_user "screen -p 0 -S $SCREEN -X eval 'stuff \"save-off\"\015'"
		as_user "screen -p 0 -S $SCREEN -X eval 'stuff \"save-all\"\015'"
#		sync
		sleep 10
	else
		echo "$SERVICE was not running. Not suspending saves."
	fi
}

mc_saveon() {
	if [ screen_running ]; then
		echo "$SERVICE is running... re-enabling saves"
		as_user "screen -p 0 -S $SCREEN -X eval 'stuff \"save-on\"\015'"
	else
		echo "$SERVICE was not running. Not resuming saves."
	fi
}

mc_stop() {
	if [ screen_running ]; then
		echo "$SERVICE is running... stopping."
		as_user "screen -p 0 -S $SCREEN -X eval 'stuff \"save-all\"\015'"
		sleep 5
		as_user "screen -p 0 -S $SCREEN -X eval 'stuff \"stop\"\015'"
		sleep 10
	else
		echo "$SERVICE was not running."
	fi
	if [ screen_running ]; then
		echo "$SERVICE is still running."
		echo "$SERVICE will now receive a second stop command."
		as_user "screen -p 0 -S $SCREEN -X eval 'stuff \"stop\"\015'"
		sleep 10
	fi
	if [ screen_running ]; then
		echo "$SERVICE could not be shut down... still running."
	else
		echo "$SERVICE is shut down."
	fi
}

mc_backup() {
# Switching this to rsync.
#	as_user "mkdir -p $WHOLEBACKUP"
#	path=`datepath $WHOLEBACKUP/mine_`
#	as_user "cp -r $MCPATH $path"
}

case "$1" in
	start)
		# Starts the server
		mc_start
		;;
	stop)
		# Stops the server
		as_user "screen -p 0 -S $SCREEN -X eval 'stuff \"say SERVER SHUTTING DOWN!\"\015'"
		mc_stop
		;;
	restart)
		# Restarts the server
		as_user "screen -p 0 -S $SCREEN -X eval 'stuff \"say SERVER REBOOT IN 5 SECONDS.\"\015'"
		mc_stop
		echo "sleeping for 5 seconds. Before starting server again."
		sleep 5
		mc_start
		;;
	backup)
		# Backups world
		as_user "screen -p 0 -S $SCREEN -X eval 'stuff \"say Backing up world.\"\015'"
		mc_saveoff
		mc_backup
		mc_saveon
		as_user "screen -p 0 -S $SCREEN -X eval 'stuff \"say Backup complete.\"\015'"
		;;
	connected)
		# Lists connected users
		as_user "screen -p 0 -S $SCREEN -X eval 'stuff list\015'"
		sleep 3s
		tac $MCPATH/server.log | grep -m 1 "Connected"
		;;
	status)
		# Shows server status
		if [ screen_running ]; then
			echo "$SERVICE is running."
		else
			echo "$SERVICE is not running."
		fi
		;;
	version)
		echo Craftbukkit version `awk '/Craftbukkit/ {sub(/\)/, ""); print $12}' $MCPATH/server.log`
		;;
	help)
		echo "Usage: /etc/init.d/minecraft command"
		echo ""
		echo "start - Starts the server"
		echo "stop - stops the server"
		echo "restart - restarts the server"
		echo "backup - backups the server up to $BACKUP"
		echo "connected - lists connected users"
		echo "status - Shows server status"
		echo "version - returs Bukkit version"
		;;
	*)
		echo "No such command see /etc/init.d/minecraft help"
		exit 1
		;;
esac

exit 0
